# paperfs.container
[Unit]
Description=paperfs.rs container

[Container]
Image=localhost/paperfs:latest
EnvironmentFile=%h/paperfs/.env
WorkingDir=/paperfs
Volume=paperfs:/paperfs
PublishPort=3000:3000

[Install]
WantedBy=default.target
